# system_files/usr/lib/systemd/user/udiskie.service
[Unit]
Description=udiskie automounter

[Service]
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target