# system_files/usr/lib/systemd/user/wl-clip-persist.service
[Unit]
Description=Persistent wayland clipboard

[Service]
ExecStart=/usr/bin/wl-clip-persist --clipboard regular

[Install]
WantedBy=default.target