# system_files/etc/profile.d/fastfetch-aliases.sh
#!/bin/bash

alias fastfetch="fastfetch --logo /usr/share/tartaria/ascii-logo/logo.txt"
alias neofetch="fastfetch --logo /usr/share/tartaria/ascii-logo/logo.txt"
alias hyfetch="fastfetch --logo /usr/share/tartaria/ascii-logo/logo.txt"